# board/fischertechnik/TXT/rootfs/etc/profile.d/ftc.sh
# ftc-specific environment variables
export PYTHONPATH=/opt/ftc
export DISPLAY=:0 
